# utils/supervisord.conf
[program:tomcat]
directory=/usr/local/tomcat
command=/usr/local/tomcat/bin/catalina.sh run
logfile_maxbytes = 50MB
logfile_backups=10
startsecs=10
